# setup/mltracker.service
[Unit]
Description=MLTracker — Gunicorn WSGI server
After=network.target __REDIS_SERVICE__
Requires=__REDIS_SERVICE__

[Service]
User=__REPO_USER__
Group=__REPO_USER__
WorkingDirectory=__REPO_DIR__/backend
EnvironmentFile=/etc/mltracker.env
ExecStart=__REPO_DIR__/backend/venv/bin/gunicorn \
    -c gunicorn.conf.py \
    "app:create_app()"
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
